# solvra_shell_v0.1.0/extras/systemd/solvra-shell-launcher.service
[Unit]
Description=Solvra Shell Launcher
After=solvra-compositor.service

[Service]
ExecStart=/usr/local/bin/solvra_shell_launcher
Restart=on-failure
Environment=RUST_LOG=info

[Install]
WantedBy=default.target
