# solvra_shell_v0.1.0/extras/systemd/solvra-shell-settings.service
[Unit]
Description=Solvra Shell Settings Panel
After=solvra-compositor.service

[Service]
ExecStart=/usr/local/bin/solvra_shell_settings
Restart=on-failure
Environment=RUST_LOG=info

[Install]
WantedBy=default.target
