# nova_shell_v0.1.0/extras/systemd/nova-shell-launcher.service
[Unit]
Description=Nova Shell Launcher
After=nova-compositor.service

[Service]
ExecStart=/usr/local/bin/nova_shell_launcher
Restart=on-failure
Environment=RUST_LOG=info

[Install]
WantedBy=default.target
