# nova_gui/extras/systemd/nova-compositor.service
[Unit]
Description=Nova GUI Compositor
After=graphical.target

[Service]
ExecStart=/usr/local/bin/nova_compositor
Restart=on-failure
Environment=RUST_LOG=info

[Install]
WantedBy=graphical.target
