# solvra_gui/extras/systemd/solvra-compositor.service
[Unit]
Description=Solvra GUI Compositor
After=graphical.target

[Service]
ExecStart=/usr/local/bin/solvra_compositor
Restart=on-failure
Environment=RUST_LOG=info

[Install]
WantedBy=graphical.target
